# crontab
PATH=/usr/local/bin:/usr/bin:/usr/local/sbin:/usr/sbin:/bin
SHELL=/bin/bash
* * * * * /home/pi/speedtest/create_speedtest_csv.sh >> /home/pi/speedtest/log.log 2>&1
